# lib/dpkg/info/php8.3-fpm.conffiles
/etc/apache2/conf-available/php8.3-fpm.conf
/etc/init.d/php8.3-fpm
/etc/logrotate.d/php8.3-fpm
/etc/php/8.3/fpm/php-fpm.conf
/etc/php/8.3/fpm/pool.d/www.conf
